# net/dnsdist/pkg/dnsdist.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/bin/dnsdist -d -u _dnsdist"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
